# deploy/com.thp.cc.cmi-plausibility.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.thp.cc.cmi-plausibility</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>tools/plausibility_check.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/thp/source/repos/wp-state-machine</string>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/thp/Library/Logs/cc/plausibility/plausibility.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/thp/Library/Logs/cc/plausibility/plausibility.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
